# deploy/crawler.service
[Unit]
Description=Korean Job Board Crawler Pipeline
After=network.target docker.service
Requires=docker.service

[Service]
Type=oneshot
User=dwk1
Group=dwk1
WorkingDirectory=/home/dwk1/koreaJobCrawl
ExecStart=/home/dwk1/koreaJobCrawl/deploy/run-pipeline.sh
StandardOutput=append:/home/dwk1/crawler-pipeline.log
StandardError=append:/home/dwk1/crawler-pipeline.log

[Install]
WantedBy=multi-user.target
